# fedora-container.service
[Unit]
Description=My first nspawn container

[Service]
ExecStart=/usr/bin/systemd-nspawn -jbD /srv/fedora 3
KillMode=process
